# dist/tabrmd.service.in
# copy this file into /etc/systemd/system
[Unit]
Description=TPM2 Access Broker and Resource Manager

[Service]
Type=dbus
BusName=com.intel.tss2.Tab
Environment=G_MESSAGES_DEBUG=all
StandardOutput=syslog
ExecStart=@SBINDIR@/tpm2-abrmd

[Install]
WantedBy=runlevel3.target
