# dist/tpm2-abrmd.service.in
# copy this file into /etc/systemd/system
[Unit]
Description=TPM2 Access Broker and Resource Manager

[Service]
Type=dbus
Restart=always
RestartSec=5
EnvironmentFile=/etc/default/tpm2-abrmd
BusName=com.intel.tss2.Tab
StandardOutput=syslog
ExecStart=@SBINDIR@/tpm2-abrmd
User=tss

[Install]
WantedBy=multi-user.target
