# dist/tpm2-abrmd.service.in
[Unit]
Description=TPM2 Access Broker and Resource Management Daemon
# This condition is needed when using the device TCTI. If the
# TCP mssim is used then the condition should be commented out.
ConditionPathExistsGlob=/dev/tpm*

[Service]
Type=dbus
Restart=always
# 69=EX_UNAVAILABLE, 71=EX_OSERR, 74=EX_IOERR
RestartPreventExitStatus=69 71 74
RestartSec=5
BusName=com.intel.tss2.Tabrmd
StandardOutput=syslog
ExecStart=@SBINDIR@/tpm2-abrmd
User=tss

[Install]
WantedBy=multi-user.target
