# dist/tpm2-abrmd.service.in
[Unit]
Description=TPM2 Access Broker and Resource Management Daemon

[Service]
Type=dbus
Restart=always
RestartSec=5
BusName=com.intel.tss2.Tabrmd
StandardOutput=syslog
ExecStart=@SBINDIR@/tpm2-abrmd
User=tss

[Install]
WantedBy=multi-user.target
